# rpm/neutron-linuxbridge-agent.upstart
description "OpenStack Neutron Linux Bridge Agent"

start on stopped rc RUNLEVEL=[2345]
stop on runlevel [S016]

respawn

exec su -s /bin/sh -c "exec /usr/bin/neutron-linuxbridge-agent --log-file /var/log/neutron/linuxbridge-agent.log --config-file /usr/share/neutron/neutron-dist.conf --config-file /etc/neutron/neutron.conf --config-file /etc/neutron/plugins/linuxbridge/linuxbridge_conf.ini" neutron
