# debian/neutron-plugin-openvswitch-agent.neutron-ovs-cleanup.service.in
[Unit]
Description=OpenStack Neutron OVS cleanup

[Service]
type=oneshot
User=neutron
Group=neutron
WorkingDirectory=/var/lib/neutron
PermissionsStartOnly=true
ExecStartPre=/bin/mkdir -p /var/lock/neutron /var/log/neutron /var/lib/neutron
ExecStartPre=/bin/chown neutron:neutron /var/lock/neutron /var/log/neutron /var/lib/neutron
ExecStart=/etc/init.d/neutron-ovs-cleanup systemd-start
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

