# debian/neutron-vpn-agent.upstart
# vim:set ft=upstart ts=2 et:
description "Neutron VPN Agent"
author "Chuck Short <zulcss@ubuntu.com>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

chdir /var/run

pre-start script
  mkdir -p /var/run/neutron
  chown neutron:root /var/run/neutron
  # Check to see if openvswitch plugin in use by checking
  # status of cleanup upstart configuration
  if status neutron-ovs-cleanup; then
    start wait-for-state WAIT_FOR=neutron-ovs-cleanup WAIT_STATE=running WAITER=neutron-vpn-agent
  fi
end script

exec start-stop-daemon --start --chuid neutron --exec /usr/bin/neutron-vpn-agent -- \
       --config-file=/etc/neutron/neutron.conf --config-file=/etc/neutron/vpn_agent.ini \
       --config-file=/etc/neutron/l3_agent.ini --config-file=/etc/neutron/fwaas_driver.ini \
       --log-file=/var/log/neutron/vpn_agent.log
